# trade/trade1.service
[Unit]
Description=Streaming

[Service]
ExecStart=/home/huykbc/trade -prod=true

[Install]
WantedBy=multi-user.target
