# payment/streaming1.service
[Unit]
Description=Streaming

[Service]
ExecStart=/home/huykbc/streaming -prod=true

[Install]
WantedBy=multi-user.target
